# .zsh/profile.d/01-basics
# vim: ft=zsh
#################
# The very basics
#################
export EDITOR=nvim
export VISUAL=nvim
export PAGER=less

###################
# Terminal settings
###################
if [[ -z "$TMUX" ]] ;then
    export TERM=xterm-256color
fi
export TERMINFO=${DOT_PATH}/share/terminfo
