# .zsh/profile.d/22-homebrew
# Homebrew
HOMEBREW_PREFIX="/opt/homebrew"

if [[ "$(uname)" == "Darwin" ]]; then
    eval `${HOMEBREW_PREFIX}/bin/brew shellenv`
fi
