# .zsh/profile.d/22-rust
#####################
# Settings for rust
#####################

. "$HOME/.cargo/env"
