# .zsh/profile.d/22-rust
#####################
# Settings for rust
#####################

rustup_env_file="$HOME/.cargo/env"
[ -f "$rustup_env_file" ] && . $rustup_env_file
