# .zsh/profile.d/30-rancher-desktop
#####################
# Settings for rancher desktop
#####################

rancher_desktop_bin_dir="~/.rd/bin"
[ -d "${rancher_desktop_bin_dir}" ] && export PATH=${rancher_desktop_bin_dir}:$PATH
